# keithley_daq6510_temp_logger.conf
[program:keithely_daq6510_temp_logger]
command=/home/srgang/Keithley\ DAQ6510\ temp\ logger/Startup_bash
environment=PYTHONUNBUFFERED=1 ; let stdout logged in real time
user=srgang
autostart=true
autorestart=true
stdout_logfile=/etc/supervisor/logs/keithely_daq6510_temp_logger_out.log
stderr_logfile=/etc/supervisor/logs/keithely_daq6510_temp_logger_err.log
