# keithley_daq6510_temp_logger.conf
[program:keithely_daq6510_temp_logger]
command=/home/srgang/Keithley\ DAQ6510\ temp\ logger/Startup_bash
environment=PYTHONUNBUFFERED=1 ; let stdout logged in real time
user=srgang
autostart=true
autorestart=true

